# secrds.service
[Unit]
Description=secrds Security Monitor Agent
Documentation=https://github.com/yourorg/secrds
After=network.target
Wants=network-online.target

[Service]
Type=simple
PIDFile=/var/run/secrds.pid
ExecStart=/usr/local/bin/secrds-agent
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s

# Ensure service stays running
RemainAfterExit=no

# Security settings
User=root
Group=root
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_NET_ADMIN CAP_BPF CAP_PERFMON CAP_SYS_RESOURCE
AmbientCapabilities=CAP_SYS_ADMIN CAP_NET_ADMIN CAP_BPF CAP_PERFMON CAP_SYS_RESOURCE
NoNewPrivileges=false
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/secrds /var/run /sys/fs/bpf /sys/kernel/tracing /sys/kernel/debug

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment (optional - config is in /etc/secrds/config.yaml)
# EnvironmentFile=-/etc/secrds/env.conf

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=secrds

[Install]
WantedBy=multi-user.target

